# home/private_Library/LaunchAgents/dotfiles.recap.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dotfiles.recap</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/recap-daily</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Hourly catchup: launchd replays a missed StartInterval after wake, so
         the first fire after a sleeping midnight still runs. The wrapper exits
         early when yesterday's recap already exists, making extra fires cheap. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/recap-daily.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/recap-daily.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
    </dict>
</dict>
</plist>
